# backend/processing/cron/event_time_update_job.crontab
0 1 * * * ec2-user bash ~/.orbsight/cron/demo_event_time_rotator.sh orb_tenant_sportsarena >> ~/.orbsight/logs/event_time_cron_sportsarena.log
