# backend/processing/cron/observer_messager_job.crontab
*/5 * * * * ec2-user bash ~/obsights/processing/cron/active_observer_publisher.sh >> ~/.obsights/observer_cron.log
