# backend/processing/cron/observer_messager_job.crontab
*/5 * * * * ec2-user bash ~/.orbsight/cron/active_observer_publisher.sh >> ~/.orbsight/observer_cron.log
